# suspend-prevent.service
[Unit]
Description=Prevent suspend on lid close when on AC

[Service]
ExecStart=/usr/local/bin/suspend-prevent --forever
